# supervisord_xuanji_loops.conf
[program:xuanji_autonomous]
directory=/workspaces/-AI-2.0
command=/bin/bash -lc 'source .venv/bin/activate && python -u autonomous_run.py'
autostart=true
autorestart=true
startretries=0
stopasgroup=true
killasgroup=true
stdout_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_autonomous.out.log
stderr_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_autonomous.err.log
environment=PYTHONUNBUFFERED="1"

[program:xuanji_predict]
directory=/workspaces/-AI-2.0
command=/bin/bash -lc 'source .venv/bin/activate && python -u ssq_predict_cycle.py'
autostart=true
autorestart=true
startretries=3
stdout_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_predict.out.log
stderr_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_predict.err.log
environment=PYTHONUNBUFFERED="1",SSQ_CONTINUOUS="1"

[program:xuanji_optimize]
directory=/workspaces/-AI-2.0
command=/bin/bash -lc 'source .venv/bin/activate && python -u auto_optimize_cycle.py'
autostart=true
autorestart=true
startretries=3
stdout_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_optimize.out.log
stderr_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_optimize.err.log
environment=PYTHONUNBUFFERED="1"

[program:xuanji_status_updater]
directory=/workspaces/-AI-2.0
; 每60秒刷新一次状态报告
command=/bin/bash -lc 'while true; do source .venv/bin/activate && python -u tools/update_status_report.py; sleep 60; done'
autostart=true
autorestart=true
startretries=3
stdout_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_status_updater.out.log
stderr_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_status_updater.err.log
environment=PYTHONUNBUFFERED="1"

[program:xuanji_evaluator]
directory=/workspaces/-AI-2.0
; 每4小时生成一次评估报告
command=/bin/bash -lc 'while true; do source .venv/bin/activate && python -u ssq_evaluate.py; sleep 14400; done'
autostart=true
autorestart=true
startretries=3
stdout_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_evaluator.out.log
stderr_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_evaluator.err.log
environment=PYTHONUNBUFFERED="1"

[program:xuanji_report]
directory=/workspaces/-AI-2.0
command=/bin/bash -lc 'source .venv/bin/activate && uvicorn celestial_nexus.report_frontend:app --host 0.0.0.0 --port 8080'
autostart=true
autorestart=true
startretries=3
stdout_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_report.out.log
stderr_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_report.err.log
environment=PYTHONUNBUFFERED="1"

[program:xuanji_static]
directory=/workspaces/-AI-2.0/static
command=/bin/bash -lc 'source /workspaces/-AI-2.0/.venv/bin/activate && python -m http.server 8088 --bind 0.0.0.0'
autostart=true
autorestart=true
startretries=3
stdout_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_static.out.log
stderr_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_static.err.log
environment=PYTHONUNBUFFERED="1"

[program:xuanji_logrotate]
directory=/workspaces/-AI-2.0
command=/bin/bash -lc 'while true; do bash tools/rotate_logs.sh; sleep 3600; done'
autostart=true
autorestart=true
startretries=3
stdout_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_logrotate.out.log
stderr_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_logrotate.err.log
environment=PYTHONUNBUFFERED="1"

[program:xuanji_public_tunnels]
directory=/workspaces/-AI-2.0
; 通过 pyngrok 暴露 8000/8080/8088 端口到公网（如已配置 NGROK_AUTH_TOKEN）
command=/bin/bash -lc 'source .venv/bin/activate && python -u tools/start_public_tunnels.py'
autostart=true
autorestart=true
startretries=3
stdout_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_public_tunnels.out.log
stderr_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_public_tunnels.err.log
environment=PYTHONUNBUFFERED="1"


[program:xuanji_person]
directory=/workspaces/-AI-2.0
command=/bin/bash -lc 'source .venv/bin/activate && python -u person_predict_cycle.py'
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stdout_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_person.out.log
stderr_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_person.err.log
environment=PYTHONUNBUFFERED="1"

[program:xuanji_internet_selftest]
directory=/workspaces/-AI-2.0
; 每10分钟执行一次互联网自检，写入 static/internet_self_test.json
command=/bin/bash -lc 'while true; do source .venv/bin/activate && python -u tools/self_test_internet.py; sleep 600; done'
autostart=true
autorestart=true
startretries=3
stdout_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_internet_selftest.out.log
stderr_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_internet_selftest.err.log
environment=PYTHONUNBUFFERED="1"

[program:xuanji_tuner]
directory=/workspaces/-AI-2.0
; 每6小时运行一次自动调优器，更新策略权重与号码先验
command=/bin/bash -lc 'while true; do source .venv/bin/activate && python -u ssq_auto_tuner.py --data ssq_history.csv --train-ai --window 2000 --ema 0.6; sleep 21600; done'
autostart=true
autorestart=true
startretries=3
stdout_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_tuner.out.log
stderr_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_tuner.err.log
environment=PYTHONUNBUFFERED="1"

[program:xuanji_tune_on_change]
directory=/workspaces/-AI-2.0
; 监控 patterns_knowledge.json 变更并立即触发一次快速调优
command=/bin/bash -lc 'source .venv/bin/activate && python -u tools/watch_patterns_and_tune.py'
autostart=true
autorestart=true
startretries=3
stdout_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_tune_on_change.out.log
stderr_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_tune_on_change.err.log
environment=PYTHONUNBUFFERED="1"

[program:xuanji_fetch_ssq]
directory=/workspaces/-AI-2.0
; 每天执行若干次（此处每3小时一次），可改为开奖日策略
command=/bin/bash -lc 'while true; do source .venv/bin/activate && python -u tools/fetch_latest_ssq.py; sleep 10800; done'
; 注意：已接入真实 TSV 导入，默认关闭 mock 拉取；如无 TSV 可临时开启
autostart=false
autorestart=true
startretries=3
stdout_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_fetch_ssq.out.log
stderr_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_fetch_ssq.err.log
environment=PYTHONUNBUFFERED="1"

[program:xuanji_import_tsv]
directory=/workspaces/-AI-2.0
; 定期从系统文件 reports/双色球 导入真实开奖：每30分钟检查一次；脚本内部按北京时间开奖窗口决定是否实际导入
command=/bin/bash -lc 'while true; do source .venv/bin/activate && python -u tools/import_ssq_from_tsv.py; sleep 1800; done'
autostart=true
autorestart=true
startretries=3
stdout_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_import_tsv.out.log
stderr_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_import_tsv.err.log
environment=PYTHONUNBUFFERED="1"

[program:xuanji_culturaldl_trainer]
directory=/workspaces/-AI-2.0
; 每12小时训练一次深度文化模型
command=/bin/bash -lc 'while true; do source .venv/bin/activate && python -u tools/train_cultural_deep.py ssq_history.csv; sleep 43200; done'
autostart=true
autorestart=true
startretries=3
stdout_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_culturaldl_trainer.out.log
stderr_logfile=/workspaces/-AI-2.0/logs/supervisor/xuanji_culturaldl_trainer.err.log
environment=PYTHONUNBUFFERED="1"
